# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
loglevel=info

[program:php-artisan-serve]
command=php artisan serve --host=0.0.0.0 --port=%(ENV_PORT)s
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=ENV_PORT="8000"

[program:queue-listen]
command=php artisan queue:listen --tries=1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:reverb]
command=php artisan reverb:start
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:vite-dev]
command=npm run dev
directory=/app
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
