# supervisord.conf
; supervisord.conf (put at repo root and copied into image)
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info

[program:php-artisan-serve]
command=php artisan serve --host=0.0.0.0 --port=8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/php-artisan-serve.log
stderr_logfile=/var/log/php-artisan-serve.err

[program:queue-listen]
command=php artisan queue:listen --tries=1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/queue-listen.log
stderr_logfile=/var/log/queue-listen.err

[program:reverb]
command=php artisan reverb:start
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/reverb.log
stderr_logfile=/var/log/reverb.err

[program:vite-dev]
command=npm run dev
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/vite-dev.log
stderr_logfile=/var/log/vite-dev.err
